# ragnarok-base/etc/cron.daily/update-manpages
#!/bin/sh

# Update mandoc's database.
/usr/sbin/makewhatis /usr/share/man/
